# others/AIS3 2024/EOF/#5 Pwn2Own/src/v15/_firmware-v15.extracted/pkg/systemd/display_manager.service
[Unit]
Description=start display manager
After=x11.service
Requires=x11.service

[Service]
Type=simple
User=eof
Group=eof
Environment=DISPLAY=:0
ExecStart=/usr/bin/display_manager
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
